# provisioning/roles/single.stucco/templates/etc/supervisor/conf.d/stucco-ui.conf
[supervisord]
loglevel=info

[program:stucco-ui]
command=node server.js
directory={{ stucco_dir }}/ui
process_name=stucco-ui
numprocs=1
autostart=true
autorestart=true
