# provisioning/roles/stucco.graph/templates/supervisor/db-monitor.conf
[supervisord]
loglevel=warn

[program:db-monitor]
command=bash /{{ stucco_dir }}/orientdb_monitor.sh
process_name=db-monitor
numprocs=1
autostart=true
autorestart=true
priority=500
startretries=10