# provisioning/roles/stucco.collectors/templates/supervisor/stucco-scheduler.conf
[supervisord]
loglevel=warn

[program:stucco-scheduler]
command={{ stucco_collector_dir }}/start.sh
directory={{ stucco_dir }}
numprocs=1
autostart=true
autorestart=true
priority=900
startretries=10
