# provisioning/roles/stucco.ui/templates/supervisor/stucco-db-rest-service.conf
[supervisord]
loglevel=warn

[program:stucco-db-rest-service]
command=java -jar target/DBConnectionService-1.0.0.jar
directory={{ stucco_dir }}/graph-db-connection-service
process_name=stucco-db-rest-service
numprocs=1
autostart=true
autorestart=true
priority=500
startretries=10
environment=STUCCO_DB_TYPE=POSTGRESQL,STUCCO_DB_CONFIG="{{ stucco_dir }}/rt/graph-db-connection/config/postgresql.yml",SITU_DB_CONFIG="{{ stucco_dir }}/rt/graph-db-connection/config/elasticsearch.yml",STUCCO_DB_INDEX_CONFIG="",STUCCO_UI_CONFIG="{{ stucco_dir }}/ui/dev_config.json"
