# provisioning/roles/jgoodall.golang/templates/etc/profile.d/golang.sh
export GOROOT={{ goroot }}
export GOBIN=$GOROOT/bin
export PATH=$GOBIN:$PATH