# provisioning/roles/ornl-sava.stucco/templates/etc/supervisor/conf.d/stucco-doc.conf
[supervisord]
loglevel=debug

[program:stucco-doc]
command={{ stucco_bin }}/doc-service -port={{ doc_service_port }}
numprocs=1
autostart=true
autorestart=true
