# provisioning/roles/ornl-sava.stucco/templates/etc/supervisor/conf.d/stucco-ui.conf
[supervisord]
loglevel=debug

[program:stucco-ui]
command=http-server -p 8000
directory={{ stucco_dir }}/ui/dist
numprocs=1
autostart=true
autorestart=true
