# provisioning/roles/ornl-sava.stucco/templates/etc/supervisor/conf.d/stucco-scheduler.conf
[supervisord]
loglevel=debug

[program:stucco-scheduler]
command={{ stucco_dir }}/collectors/scheduler-vm.sh production
directory={{ stucco_dir }}
numprocs=1
autostart=true
autorestart=true
