# provisioning/roles/stucco.doc/templates/supervisor/stucco-doc.conf
[supervisord]
loglevel=info

[program:stucco-doc]
command={{ stucco_bin }}/doc-service -port={{ stucco_doc_port }}
directory={{ stucco_dir }}
numprocs=1
autostart=true
autorestart=true
