# provisioning/roles/ornl-sava.titan/templates/etc/init/titan.conf
description "titan graph database"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  {{ titan_path }}/bin/titan.sh -v start
end script
